# recipes-bsp/asus-overlay/files/imx8mq-pv100a/lib/systemd/system/rtcsync.service
[Unit]
Description=Sync system to MCU HW clock
After=ntpsync.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/hwclock_mcu -w

[Install]
WantedBy=multi-user.target
