# recipes-tpm/tpm2-tss/files/tpm2-tss.service
[Unit]
Description=ttyGPS Daemon
After=network.target

[Service]
ExecStart=/etc/init.d/tpm2-tss.sh start

[Install]
WantedBy=multi-user.target
