# recipes-bsp/asus-overlay/files/common/etc/udev/rules.d/99-asus-lte-manager.rules
KERNEL!="cdc-wdm*", GOTO="asus_lte_manager_end"

ACTION=="add", RUN+="/sbin/lte.sh start"

ACTION=="remove", RUN+="/sbin/lte.sh stop"

LABEL="asus_lte_manager_end"
