# recipes-bsp/asus-overlay/files/common/lib/systemd/system/fotaclient.service
[Unit]
Description=Fota Client Service
DefaultDependencies=no
After=resize-helper.service

[Service]
Type=oneshot
ExecStart=/etc/init.d/fotaclient start

[Install]
WantedBy=paths.target
